# conf/drbdmanaged.socket
[Unit]
Description=DRBDManage Service

[Socket]
ListenStream=6996

[Install]
WantedBy=default.target
