# conf/drbdmanaged.service
[Unit]
Description=DRBDManage Service

[Service]
ExecStart=/usr/bin/dbus-drbdmanaged-service
User=root

PrivateTmp=yes
